# rainfall.service
[Unit]
Description=Sprinkler System
After=network.target

[Service]
ExecStart=/opt/rainfall/server.py
Restart=always

[Install]
WantedBy=multi-user.target
